# tg-assistant/systemd/tg-prune-history.timer
[Unit]
Description=Periodic prune of tg-assistant historical messages

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
RandomizedDelaySec=10min
Persistent=true
Unit=tg-prune-history.service

[Install]
WantedBy=timers.target
